# wrapper_server.service
[Unit]
Description=Vultra VPN Server - UDP proxy for WireGuard clients
Documentation=file:///opt/vultra-vpn-server/README.md
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=vultra-vpn
Group=vultra-vpn
RuntimeDirectory=vultra-vpn-server
RuntimeDirectoryMode=0750

# Load config from env file (create from vultra-vpn-server.env.example)
EnvironmentFile=-/etc/default/vultra-vpn-server

# Logging: always enable and write to RuntimeDirectory (writable by service).
# Override WRAPPER_LOG_FILE in env file for a different path (e.g. /var/log).
Environment=WRAPPER_LOG=1
Environment=WRAPPER_LOG_FILE=/run/vultra-vpn-server/server.log

# All args from EnvironmentFile. Set PREMIUM_VERIFY_URL for premium/expiry check (wrapper calls API).
ExecStart=/opt/vultra-vpn-server/wrapper_server ${WRAPPER_BASE_PORT} ${WRAPPER_NUM_PORTS} ${WG_SERVER} ${WG_PORT} ${WRAPPER_KEY_HEX}

Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vultra-vpn-server

# Bind to privileged port (e.g. 78) as vultra-vpn
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadOnlyPaths=/etc/default
ReadWritePaths=/opt/vultra-vpn-server
InaccessibleDirectories=/root /home
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictAddressFamilies=AF_INET AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
PrivateMounts=yes

[Install]
WantedBy=multi-user.target
